# ha_notification_tap/rootfs/etc/services.d/notify-tap/finish
#!/usr/bin/execlineb -S1
foreground { s6-svscanctl -t /var/run/s6/services }
