# ha_notification_tap/rootfs/etc/services.d/notify-tap/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the notify-tap service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Redirect stdout/stderr to Home Assistant's logging system
exec 2>&1

bashio::log.info "Starting Notification Tap service..."
cd /usr/src/app
exec python3 run.py
