# ha_notification_tap/rootfs/etc/services.d/notify-tap/run
#!/usr/bin/with-contenv bashio

cd /usr/src/app || bashio::exit.nok "Failed to change directory"

exec python3 -u run.py
