# ha_notification_tap/rootfs/etc/services.d/notify-tap/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the notify-tap service
# ==============================================================================

export PYTHONUNBUFFERED=1
cd /usr/src/app
exec python3 run.py
