# ha_notification_tap/rootfs/etc/services.d/notify-tap/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the notify-tap service
# ==============================================================================

# Ensure output is not buffered
export PYTHONUNBUFFERED=1

# Redirect all output to stdout for s6-overlay
exec 2>&1

# Use bashio for logging
bashio::log.info "Starting Notification Tap service..."

# Set up error handling
trap 'bashio::exit.nok "Unexpected error occurred"' ERR

# Ensure we're in the right directory
cd /usr/src/app || bashio::exit.nok "Could not change directory"

# Start the service with unbuffered output
exec python3 -u run.py
