# ha_notification_tap/rootfs/etc/services.d/notify-tap/run
#!/usr/bin/with-contenv bashio

# Enable debug if configured
if bashio::config.true 'debug'; then
    bashio::log.level "debug"
fi

cd /usr/src/app || bashio::exit.nok "Could not change directory"

bashio::log.info "Starting Notification Tap service..."
exec python3 -u run.py
