# ha_notification_tap/rootfs/etc/services.d/notify-tap/run
#!/usr/bin/with-contenv bashio

# Set up proper output handling
exec 1> >(exec logger -s -t $(basename $0)) 2>&1

cd /usr/src/app || exit 1

# Run with unbuffered output
PYTHONUNBUFFERED=1 exec python3 run.py
