# ha_notification_tap/rootfs/etc/services.d/notify-tap/run
#!/usr/bin/with-contenv bashio

cd /usr/src/app || exit 1

# Redirect stderr to stdout
exec 2>&1

# Run with systemd-compatible output
exec python3 -u run.py
