# ha_notification_tap/rootfs/etc/services.d/notify-tap/run
#!/usr/bin/with-contenv bashio

# Enable debug if configured
if bashio::config.true 'debug'; then
    bashio::log.level "debug"
fi

cd /usr/src/app || exit 1

bashio::log.info "Starting Notification Tap service..."
PYTHONUNBUFFERED=1 exec python3 -u run.py 2>&1
